# eproba/crontab
0 0 * * * /home/app/web/autobackup.sh >> /var/log/cron/backups.log 2>&1
